# COPY/bootstrap.conf
[program:bootstrap]
command=/bin/bash /usr/local/bin/bootstrap.sh
priority=1
startsec=0
autorestart=false